# rasp/sensor_reader-startup.service
[Unit]
Description=Start Sensor Reader Script

[Service]
ExecStart=/home/pi/proyecto/startup.sh

[Install]
WantedBy=multi-user.target
